# system/halium/bluetooth-touch-tenshi.conf
description "Bluetooth initialization for BQ Aquaris U Plus (tenshi)"

start on starting bluetooth

expect fork
exec /bin/sh /system/usr/share/bluetooth-touch/tenshi
